# templates/app.service.erb
[Unit]
Description=<%= app_name %> application server
After=network.target

[Service]
Type=simple
User=<%= app_user %>
WorkingDirectory=/var/www/<%= app_name %>/current
Environment=RAILS_ENV=production
Environment=PATH=/home/<%= app_user %>/.rbenv/shims:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/<%= app_user %>/.rbenv/shims/bundle exec puma -C config/puma.rb
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
